# assets/systemd/firstmate-autostart.service
# firstmate-autostart.service - TEMPLATE, not an installed unit.
#
# Boot brings up `herdr server`, which is headless and materialises no panes, so
# nothing was actually supervising the fleet until a human ran `herdr` and
# attached. This unit runs the missing client-side step over herdr's socket API.
#
# __FM_ROOT__ is a placeholder, NOT a shell variable: replace it with the
# absolute path of the firstmate home before installing. The install, rollback,
# and verification commands live in docs/firstmate-autostart.md - use those
# rather than copying this file by hand.

[Unit]
Description=Start the firstmate primary agent in the herdr server
Documentation=file://__FM_ROOT__/docs/firstmate-autostart.md
# Wants, not Requires: this is a best-effort follower. Requires would also tear
# this unit down whenever herdr-server restarts, and a herdr-server restart is
# precisely when a resurrected-but-flagless firstmate is already handled
# elsewhere (see docs/claude-resume-shim.md). After= orders us behind the server
# PROCESS; the readiness poll in ExecStart handles the socket not yet answering.
Wants=herdr-server.service
After=herdr-server.service
# The agent registers with Anthropic's remote-control service at launch, so it
# must not start before the network is usable. Know the limit of these two
# lines: they order against the USER manager's network-online.target, and a
# user manager typically has NO provider for it (verified 2026-07-20 on the
# captain's host: no unit file, permanently inactive), so on such hosts they
# are inert - a user unit cannot order against system targets. The wait that
# actually enforces "network first" is fm-autostart.sh's own bounded
# reachability gate (exit 5); these lines engage only where a user-scope
# provider exists.
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
# systemd user units get a minimal PATH, and herdr lives in ~/.local/bin, so
# without this line the script dies at boot with "herdr not found on PATH"
# (data/learnings.md, 2026-07-20). %h is the user's home directory, expanded
# by systemd, so the template stays host-neutral.
Environment="PATH=%h/.local/bin:/usr/local/bin:/usr/bin:/bin"
ExecStart=__FM_ROOT__/bin/fm-autostart.sh --timeout 120
# Must exceed the script's own budget (--timeout 120 for the server wait, plus
# --net-timeout 120 for the network gate, plus its 20s start confirmation) so
# systemd never kills the wait before it can report why.
TimeoutStartSec=420

[Install]
WantedBy=default.target
